# quadlets/drogon-postgres.container
# containers/drogon-postgres.container
# Podman quadlet definition for the PostgreSQL database
# Install to /etc/containers/systemd/drogon-postgres.container
#
# Equivalent to:
#   podman run --name drogon-postgres --rm \
#     -e POSTGRES_PASSWORD=secret \
#     -e POSTGRES_DB=vote \
#     -v /home/peschke/code/private/Vote-Backend/sql:/docker-entrypoint-initdb.d \
#     -p 5432:5432 \
#     -d postgres:18

[Container]
ContainerName=drogon-postgres
Image=postgres:18
Environment=POSTGRES_PASSWORD=secret
Environment=POSTGRES_DB=vote
Volume=%h/code/private/vote-backend/sql:/docker-entrypoint-initdb.d:Z
Network=vote.network
PublishPort=5432:5432

[Service]
Restart=on-failure
# Wait until PostgreSQL is truly ready to accept connections
# before signaling the service as started.
ExecStartPost=/bin/sh -c 'until podman exec drogon-postgres pg_isready -h localhost -p 5432 -U postgres; do sleep 1; done'

[Install]
