# quadlets/vote-backend.container
# containers/vote-backend.container
# Podman quadlet definition for the vote-backend application
# Install to /etc/containers/systemd/vote-backend.container
#
# Build the image first:
#   podman build -f Containerfile -t vote-backend .
#
# The database must be running before this container starts.
# Config is loaded from /etc/vote/config.json inside the container.

[Container]
ContainerName=vote-backend
Image=localhost/vote-backend:0.1
Network=vote.network
PublishPort=8848:8848

[Service]
Restart=on-failure

[Unit]
After=drogon-postgres.service
Requires=drogon-postgres.service

[Install]
