# ansible/roles/deploy-local/templates/drogon-postgres.container.j2
# containers/drogon-postgres.container
# Podman quadlet definition for the PostgreSQL database.
# Rendered by the deploy-local role from this template.
#
# When `vote_backend_postgres_persistent` is true the PostgreSQL data directory
# is bind-mounted from {{ vote_backend_postgres_data_dir }} so that the database
# survives container restarts. Otherwise the database is ephemeral (the data
# lives only inside the container and is reset on every (re)deploy).
#
# PostgreSQL >= 18 requires the major version ("18") to appear in the data
# directory path (see docker-library/postgres#1259). This applies to the
# in-container PGDATA, so PGDATA is set to /var/lib/postgresql/data/18 and the
# host bind-mount (when persistent) maps to that same path. Both therefore
# include "/18".
#
# Equivalent to (ephemeral):
#   podman run --name drogon-postgres -e POSTGRES_PASSWORD=secret \
#     -e POSTGRES_DB=vote -e PGDATA=/var/lib/postgresql/data/18 \
#     -v ~/.local/share/vote-backend/sql:/docker-entrypoint-initdb.d \
#     -p 5432:5432 -d postgres:18
#
# Generated systemd unit: drogon-postgres.service

[Container]
ContainerName=drogon-postgres
Image=postgres:18
Environment=POSTGRES_PASSWORD=secret
Environment=POSTGRES_DB=vote
Environment=PGDATA=/var/lib/postgresql/data/18
Volume=%h/.local/share/vote-backend/sql:/docker-entrypoint-initdb.d:Z
{% if vote_backend_postgres_persistent %}
Volume=%h/.local/share/vote-backend/postgres/18:/var/lib/postgresql/data/18:Z
{% endif %}
Network=vote.network
PublishPort=5432:5432

[Service]
Restart=on-failure
# Wait until PostgreSQL is truly ready to accept connections
# before signaling the service as started.
ExecStartPost=/bin/sh -c 'until podman exec drogon-postgres pg_isready -h localhost -p 5432 -U postgres; do sleep 1; done'

[Install]
